# etc/insserv/overrides/memlockd
### BEGIN INIT INFO
# Provides:          memlockd
# Required-Start:    $remote_fs $syslog
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      
### END INIT INFO
