# debian/wiperam.postinst
#!/bin/sh

set -e

if [ -x "/etc/init.d/kexec-load" ]; then
   update-rc.d kexec-load defaults
fi

#DEBHELPER#
