# debian/wiperam.postrm
#!/bin/sh

set -e

update-rc.d kexec-load defaults || true

#DEBHELPER#
